# A9/移植部分/gprs短信模块/gprs测试程序_中文/gprs+fs/rootfs/etc/udev/rules.d/10-add-sd.rules
ACTION!="add",GOTO="farsight"
KERNEL=="mmcblk[0-9]p[0-9]",RUN+="/sbin/mount-sd.sh %k"
LABEL="farsight"
