# A9/移植部分/gprs短信模块/gprs测试程序_中文/gprs+fs/rootfs/etc/udev/rules.d/11-add-usb.rules
ACTION!="add",GOTO="farsight"
KERNEL=="sd[a-z][0-9]",RUN+="/sbin/mount-usb.sh %k"
LABEL="farsight"
