# A9/移植部分/gprs短信模块/gprs测试程序_中文/gprs+fs/rootfs/etc/udev/rules.d/11-remove-usb.rules
ACTION !="remove",GOTO="farsight"
SUBSYSTEM!="block",GOTO="farsight"
KERNEL=="sd[a-z][0-9]",RUN+="/sbin/umount-usb.sh"
LABEL="farsight"

